# btpan.service
[Unit]
After=bluetooth.service
PartOf=bluetooth.service

[Service]
ExecStart=/usr/local/bin/btpan

[Install]
WantedBy=bluetooth.target
